# deploy/attendx.service
[Unit]
Description=ATTENDX Flask Backend
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=attendx
Group=attendx
WorkingDirectory=/home/attendx/attendx/backend
Environment="PATH=/home/attendx/attendx/backend/venv/bin:/usr/bin"
EnvironmentFile=/home/attendx/attendx/backend/.env
ExecStart=/home/attendx/attendx/backend/venv/bin/gunicorn \
    --bind 127.0.0.1:5001 \
    --workers 2 \
    --timeout 120 \
    --access-logfile /home/attendx/logs/access.log \
    --error-logfile /home/attendx/logs/error.log \
    "run:app"
Restart=always
RestartSec=5

# Create log directory
ExecStartPre=/bin/mkdir -p /home/attendx/logs

[Install]
WantedBy=multi-user.target
